# docker/supervisord-worker.conf
[unix_http_server]
file=/tmp/supervisor.sock
username = dummy
password = dummy

[inet_http_server]
port=127.0.0.1:9001
username = dummy
password = dummy

[supervisord]
nodaemon=true
loglevel=warn

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username = dummy
password = dummy


# [eventlistener:flush_signalbus_trigger]
# command=%(ENV_APP_ROOT_DIR)s/trigger_supervisor_process.py signalbus_flush 60
# directory=%(ENV_APP_ROOT_DIR)s
# events=TICK_5

# [program:signalbus_flush]
# command=flask signalbus flush -w 10
# directory=%(ENV_APP_ROOT_DIR)s
# autostart=false
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes = 0
# redirect_stderr=true
# startsecs=0
# exitcodes=0
# autorestart=false


[program:flush_all]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_all
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:consume_messages]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh consume_messages
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:consume_chore_messages]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh consume_chore_messages
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:roll_worker_turns]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh roll_worker_turns
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:fetch_debtor_infos]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh fetch_debtor_infos
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:process_pristine_collectors]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh process_pristine_collectors
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_debtor_info_documents]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_debtor_info_documents
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_debtor_locator_claims]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_debtor_locator_claims
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_trading_policies]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_trading_policies
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_worker_accounts]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_worker_accounts
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_interest_rate_changes]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_interest_rate_changes
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_needed_worker_accounts]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_needed_worker_accounts
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_recently_needed_collectors]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_recently_needed_collectors
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_account_locks]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_account_locks
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_creditor_participations]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_creditor_participations
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_dispatching_statuses]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_dispatching_statuses
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_worker_collectings]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_worker_collectings
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_worker_sendings]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_worker_sendings
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_worker_receivings]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_worker_receivings
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_worker_dispatchings]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_worker_dispatchings
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_transfer_attempts]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_transfer_attempts
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000
